# src/linux/server/systemd/netremote-server.service.in
[Unit]
Description=Netremote Daemon
Requires=hostapd.service

[Service]
Type=simple
Restart=on-failure
ExecStart=${CMAKE_INSTALL_FULL_SBINDIR}/netremote-server -d

[Install]
WantedBy=default.target