# src/linux/server/systemd/netremote-server.service.in
[Unit]
Description=Netremote Daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=on-failure
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/netremote-server -d ${NETREMOTE_SERVER_XTRA_ARGS}

[Install]
WantedBy=default.target