# src/linux/external/hostap/systemd/hostapd-genconfig@.service.in
[Unit]
Description=Generate default hostapd configuration file
AssertPathExists=@CMAKE_INSTALL_FULL_SYSCONFDIR@/hostapd/hostapd-example.conf

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/cp @CMAKE_INSTALL_FULL_SYSCONFDIR@/hostapd/hostapd-example.conf @CMAKE_INSTALL_FULL_SYSCONFDIR@/hostapd/hostapd-%i.conf
ExecStart=/usr/bin/sed -i '/^#/d' @CMAKE_INSTALL_FULL_SYSCONFDIR@/hostapd/hostapd-%i.conf
ExecStart=/usr/bin/sed -i 's/ssid=/ssid=%l-ap-%i/' @CMAKE_INSTALL_FULL_SYSCONFDIR@/hostapd/hostapd-%i.conf
