# src/linux/external/hostap/systemd/hostapd.service.in
[Unit]
Description=Hostapd Daemon

[Service]
Type=forking
Restart=on-failure
ExecStart=${HOSTAP_BINARY_DIR}/hostapd -B -P /run/hostapd.pid -i wlan0 /etc/hostapd/hostapd.conf
PIDFile=/run/hostapd.pid

[Install]
WantedBy=default.target