# src/linux/external/hostap/systemd/hostapd@.service.in
[Unit]
Description=Hostapd Daemon

[Service]
Type=forked
Restart=on-failure
ExecStart=${HOSTAP_BINARY_DIR}/hostapd -B ${HOSTAP_BINARY_DIR}/${HOSTAPD_PREFIX}/.config

[Install]
WantedBy=default.target