# src/linux/external/hostap/systemd/hostapd@.service.in
[Unit]
Description=Hostapd Daemon

[Service]
Type=simple
Restart=on-failure
ExecStart=${HOSTAP_BINARY_DIR}/hostapd -B etc/hostapd/hostapd.conf

[Install]
WantedBy=default.target