# src/linux/external/hostap/systemd/hostapd.service.in
[Unit]
Description=Hostapd Daemon

[Service]
Type=forking
Restart=on-failure
ExecStart=${HOSTAP_BINARY_DIR}/hostapd -B -P /${CMAKE_INSTALL_RUNSTATEDIR}/hostapd/hostapd.pid -i wlan0 /${CMAKE_INSTALL_SYSCONFDIR}/hostapd/hostapd.conf
PIDFile=/${CMAKE_INSTALL_RUNSTATEDIR}/hostapd/hostapd.pid

[Install]
WantedBy=default.target