# src/linux/external/hostap/systemd/hostapd@.service.in
[Unit]
Description=Hostapd Daemon
Requires=sys-subsystem-net-devices-%i.device hostapd-genconfig@%i.service
After=sys-subsystem-net-devices-%i.device hostapd-genconfig@%i.service
AssertPathExists=@CMAKE_INSTALL_FULL_SYSCONFDIR@/hostapd/hostapd-%i.conf

[Service]
Type=simple
Restart=on-failure
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/hostapd -i %i @CMAKE_INSTALL_FULL_SYSCONFDIR@/hostapd/hostapd-%i.conf

[Install]
WantedBy=default.target
